# backup/crontab
0 5 * * * /app/backup.php
1 0 1 * * find /app/data/ -type f -not -name "*5.opml.gz" -mtime +30 -exec rm {} \;
1 1 1 * * find /app/data/ -type f -not -name "*05.opml.gz" -mtime +180 -exec rm {} \;
